# container-builds/ansible-navigator/configs/profile.d/decrypt.sh
#!/bin/bash

# place holder for now until I can figure out how to do this securely
